# tooling/agents/scheduler/linux/units/pella-edge-decay.service
[Unit]
Description=Pella edge-decay watchdog (rolling 90d Sharpe scan)
Documentation=file://__PELLA_LAB_ROOT__/NT8Bridge/scheduler/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PELLA_LAB_ROOT__
ExecStart=__PELLA_PYTHON__ __PELLA_LAB_ROOT__/NT8Bridge/tools/edge_decay_watchdog.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
