# tooling/agents/scheduler/linux/units/pella-discovery.service
[Unit]
Description=Pella discovery agent (nightly idea-queue worker)
Documentation=file://__PELLA_LAB_ROOT__/NT8Bridge/scheduler/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PELLA_LAB_ROOT__
ExecStart=__PELLA_PYTHON__ __PELLA_LAB_ROOT__/NT8Bridge/tools/discovery_agent.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2700
